# woof-code/rootfs-skeleton/etc/init.d/gdomap
#! /bin/sh
### BEGIN INIT INFO
# Provides:          gdomap
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start the GNUstep distributed object mapper
### END INIT INFO

case $2 in
[0-9]*) safeBOOT=$2;;
""|*) :;;
esac

PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

test -x /usr/share/GNUstep/Makefiles/GNUstep.sh || exit 0
. /usr/share/GNUstep/Makefiles/GNUstep.sh
echo "sourced /usr/share/GNUstep/Makefiles/GNUstep.sh"

DAEMON=/usr/bin/gdomap
NAME=gdomap
DESC="GNUstep distributed object mapper"

PIDFILE=/var/run/$NAME.pid

GDOMAP_ARGS=
if [ -f /etc/default/gdomap ]; then
    . /etc/default/gdomap
fi

test -x $DAEMON || exit 0
echo "$DAEMON executable"

#set -e

case "$1" in
  start)
	echo -n "Starting $DESC: "
	start-stop-daemon --start --quiet --exec $DAEMON \
		-- -I $PIDFILE $GDOMAP_ARGS
	echo "$NAME."
	;;
  stop)
	echo -n "Stopping $DESC: "
	start-stop-daemon --stop --quiet --exec $DAEMON \
		-- -I $PIDFILE $GDOMAP_ARGS
	rm -f $PIDFILE
	echo "$NAME."
	;;
  restart|force-reload)
	echo -n "Restarting $DESC: "
	start-stop-daemon --stop --quiet --exec $DAEMON \
		-- -I $PIDFILE $GDOMAP_ARGS
	rm -f $PIDFILE
	sleep 1
	start-stop-daemon --start --quiet --exec $DAEMON \
		-- -I $PIDFILE $GDOMAP_ARGS
	echo "$NAME."
	;;
  *)
	N=/etc/init.d/$NAME
	echo "Usage: $N {start|stop|restart|force-reload}" >&2
	exit 1
	;;
esac

exit 0
