# woof-code/rootfs-skeleton/lib/udev/rules.d/11-MF190-zte-modem.rules
ACTION=="add", SUBSYSTEM=="usb", ATTRS{idVendor}=="19d2", ATTRS{idProduct}=="2000", RUN+="/usr/sbin/usb_modeswitch '%b/%k'"
