# woof-code/rootfs-skeleton/etc/profile.d/30_pkgconfig
# set the PKG_CONFIG_PATH variable
#

QUIET=-q
export PKG_CONFIG_PATH=""

for x in {/usr,/opt}{,/*}/{share,lib?*,lib}/pkgconfig $HOME/lib{?*,}/pkgconfig
do

#[ -d $x ] && PKG_CONFIG_PATH="$PKG_CONFIG_PATH:$x"

[ -d "$x" ] || continue
[ -L "$x" ] && continue
[ -L "${x%/*}" ] && continue

[ "$PKG_CONFIG_PATH" ] && PKG_CONFIG_PATH="$PKG_CONFIG_PATH:$x" || PKG_CONFIG_PATH="$x"

done

export PKG_CONFIG_PATH ##+++2012-01-13
