# woof-code/rootfs-skeleton/etc/profile.d/imake
#!/bin/sh

  if [ -f /usr/lib/X11/config/Imake.tmpl ];then
  IMAKEINCLUDE="-I/usr/lib/X11/config"
elif [ -f /usr/X11R7/lib/X11/config/Imake.tmpl ];then
IMAKEINCLUDE="-I/usr/X11R7/lib/X11/config"
elif [ -f /usr/local/lib/X11/config/Imake.tmpl ];then
IMAKEINCLUDE="-I/usr/local/lib/X11/config"
  fi
export IMAKEINCLUDE
