# woof-code/rootfs-skeleton/etc/profile.d/libraries
#!/bin/ash

LD_LIBRARY_PATH='';
[ "$LD_LIBRARY_PATH" ] || LD_LIBRARY_PATH=/lib

#DIRS=`find /usr -maxdepth 5 -type d -name lib`

for n in `seq 1 1 4` ;
do
DIRS_USR=`find /usr -mindepth $n -maxdepth $n -type d -name lib | sort -d`
for dir in $DIRS_USR; do LD_LIBRARY_PATH="${LD_LIBRARY_PATH}:$dir" ; done
done

for n in `seq 1 1 4` ;
do
DIRS_OPT=`find /opt -mindepth $n -maxdepth $n -type d -name lib | sort -d`
for dir in $DIRS_OPT; do LD_LIBRARY_PATH="${LD_LIBRARY_PATH}:$dir" ; done
done

LD_LIBRARY_PATH=/usr/X11/lib:$LD_LIBRARY_PATH
export LD_LIBRARY_PATH
