# woof-code/rootfs-skeleton/etc/profile.d/libraries
#!/bin/ash

LD_LIBRARY_PATH='';
[ "$LD_LIBRARY_PATH" ] || LD_LIBRARY_PATH=/lib

#DIRS=`find /usr -maxdepth 5 -type d -name lib`
#2016-02-08: This find code is very slow on K7 1045MHz
# GNU find: 5-8 sec., bb find 20 sec.
for n in 1 2 3 4 ;
do
DIRS_USR=`/bin/find /usr -mindepth $n -maxdepth $n -type d -name lib | sort -d`
for dir in $DIRS_USR; do [ -d "$dir" ] || continue; LD_LIBRARY_PATH="${LD_LIBRARY_PATH}:$dir" ; done
done

for n in 1 2 3 4 ;
do
DIRS_OPT=`/bin/find /opt -mindepth $n -maxdepth $n -type d -name lib | sort -d`
for dir in $DIRS_OPT; do [ -d "$dir" ] || continue; LD_LIBRARY_PATH="${LD_LIBRARY_PATH}:$dir" ; done
done

LD_LIBRARY_PATH=/usr/X11/lib:$LD_LIBRARY_PATH
export LD_LIBRARY_PATH
