# woof-code/rootfs-skeleton/etc/profile.d/proc_cmdline
#!/bin/sh

[ "$DEBUG" ] || DEBUG=`grep -i 'debug' /proc/cmdline`

kernel_parameters_func(){
OUT=/dev/null;ERR=$OUT
for item1 in `cat /proc/cmdline`;
do
case $item1 in
pfix=*)
   eval $item1;export $item1
   ITEMS=${item1#*=};ITEMS=${ITEMS//,/ }
   for item2 in $ITEMS;
   do
   case $item2 in
   nodrivers:*)     :;;
   noalsa:*)        :;;
   alsa:*)          :;;
   soundcards:*)    :;;
   userdrivers:*)   :;;
   kerneldrivers:*) :;;
   alldrivers:*)    :;;
   norcsysinit)     :;;
   swaps:*) W_SWAPS=${item2#*:};W_SWAPS=${W_SWAPS//:/ };W_SWAPS=${W_SWAPS//\/dev\//};         :;;
   noswap) NO_SWAP=1;          :;;
   verbose) VERBOSE=1;VERB=-v;L_VERB=--verbose;A_VERB=-verbose;OUT=/dev/stdout;ERR=/dev/stderr;;
   debug) DEBUG=1;set -x;;
   nofb:*)          :;;
   nofbcon) NO_FBCON=1;         :;;
   fbconfirst) FBCON_FIRST=1;;
   nox)        echo 'nox' > /tmp/X.status;;
   xorgwizard) mv /etc/X11/xorg.conf /etc/X11/xorg.conf-"`date`";;

   esac;done
  ;;
*) [ "`echo "$item1" | grep '='`" ] || item1="${item1}=1"
item=${item1//\./_};
kp=${item%=*};[ "$DEBUG" ] && echo -n "${kp},"
eval $item
export $kp
;;
esac;done
[ "$DEBUG" ] && echo
}

kernel_parameters_func
