# woof-code/rootfs-skeleton/etc/profile.d/quisp
export SHSQL_DB='/root/Web-Server/Quisp'
