# woof-code/rootfs-skeleton/etc/inittab
::ctrlaltdel:/sbin/reboot

::sysinit:/etc/rc.d/rc.sysinit

:3:once:/etc/rc.d/rc.launchxwin
tty1::respawn:/etc/rc.d/rc.autologinroot
tty2::respawn:/sbin/getty -n -l /bin/autologinroot 38400 tty2
tty3::respawn:/sbin/getty 38400 tty3
tty4::respawn:/sbin/getty 38400 tty4


vc/8::respawn:/etc/init/vc8
vc/9::respawn:/etc/init/vc9




