# woof-code/rootfs-skeleton/etc/inittab
::sysinit:/etc/rc.d/rc.sysinit
tty1::respawn:/sbin/getty -n -l /bin/autologinroot 38400 tty1
tty2::respawn:/sbin/getty -n -l autologinroot 38400 tty2
tty3::respawn:/sbin/getty -n -l autologinroot 38400 tty3
::ctrlaltdel:/sbin/reboot
##tty5::respawn:/sbin/getty -n -l "login -f root" 38400 tty5
tty5::respawn:/sbin/getty -n -l autologinroot 38400 tty5
tty6::respawn:/sbin/getty -n -l autologinroot 38400 tty6
tty7::respawn:/sbin/getty -n -l autologinroot 38400 tty7
